# system/dhclient-wlan0_sta.service
[Unit]
Description=DHCP client for wlan0_sta (station interface)
After=wpa_supplicant@wlan0_sta.service
Requires=wpa_supplicant@wlan0_sta.service

[Service]
Type=oneshot
ExecStart=/sbin/dhclient -v -nw wlan0_sta
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
