# system/fre_portal.service
[Unit]
Description=FRE Node Portal (dashboard on port 80)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=aliasvava
Group=aliasvava
WorkingDirectory=/home/aliasvava/fre_node
# Adapter le chemin du venv et du code si nécessaire
ExecStart=/bin/bash -lc 'cd /home/aliasvava/fre_node && /home/aliasvava/fre_node/venv/bin/uvicorn dashboard.dashboard:app --host 0.0.0.0 --port 80'
Restart=always
Environment="PYTHONUNBUFFERED=1"
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
